# fancontrolcpp.service
[Unit]
Description=Fancontrol CPP

[Service]
Type=notify
ExecStart=/usr/sbin/fancontrolcpp

[Install]
WantedBy=multi-user.target
